# ansible/templates/phoenix.service.j2
[Unit]
Description="Phoenix Elixir Server"

[Service]
Type=simple
WorkingDirectory={{  backend_path }}
ExecStart={{ phoenix_cmd }}
{# ExecStop= #}
User={{ phoenix_user }}
Group={{ phoenix_user }}
Restart=always
StandardOutput=syslog
SyslogIdentifier=front

[Install]
WantedBy=default.target
